# rootfs/etc/services.d/crond/run
#!/usr/bin/with-contenv bash

echo "starting crond"

crond -l 2 -f

