# rootfs/etc/services.d/crond/run
#!/usr/bin/with-contenv bash

echo "starting crond"

su -c "crond -l 8 -f" root