# rootfs/etc/services.d/theia/finish
#!/usr/bin/execlineb -S0

echo "theia exiting"
s6-svscanctl -t /var/run/s6/services