# rootfs/etc/services.d/theia/run
#!/usr/bin/with-contenv bash

echo "starting theia"

su -c "node /srv/theia/app/src-gen/backend/main.js /srv/urlwatch/config --port=8081 --hostname=0.0.0.0" mediadepot